# boards/tsp-s/overlay/etc/systemd/system/pocketforge-memtest-container.service
[Unit]
Description=PocketForge memtest container (crun; proves MemoryMax enforcement via OOM)
Documentation=https://github.com/pocketforge-os/image
After=local-fs.target
AssertPathExists=/usr/lib/pocketforge/containers/memtest/config.json

[Service]
Type=simple
Slice=pocketforge-apps.slice

# Same 64M cap as hello. The memtest container deliberately allocates past it so
# the kernel OOM-kills it — proving the limit is ENFORCED, not just set. Started
# on demand by container-memtest.sh; NOT enabled at boot (it is meant to die).
# NO MemoryHigh here (unlike hello): MemoryHigh throttles the allocator into slow
# reclaim near the cap, making the OOM take tens of seconds + unpredictable
# (verified on tsp-vuo.4). A test fixture wants a fast, deterministic hard-OOM at
# MemoryMax, so we omit the soft limit.
MemoryAccounting=yes
MemoryMax=64M

ExecStartPre=-/usr/bin/crun --cgroup-manager=disabled delete --force pf-memtest
ExecStart=/usr/bin/crun --cgroup-manager=disabled run --bundle /usr/lib/pocketforge/containers/memtest pf-memtest
ExecStopPost=-/usr/bin/crun --cgroup-manager=disabled delete --force pf-memtest
Restart=no
