# boards/tsp-s/overlay/etc/systemd/system/pocketforge-watchdog-keepalive.service
[Unit]
Description=PocketForge A523 hardware watchdog keepalive (sunxi-wdt 2050000, tsp-vuo.4)
Documentation=https://github.com/pocketforge-os/image
# The sunxi-wdt auto-arms at probe and cannot be stopped; systemd's native
# RuntimeWatchdog did NOT engage on the A523 (works on the A133), so we feed the
# device directly. Start as early as possible — before the ~16s reset window.
DefaultDependencies=no
After=systemd-tmpfiles-setup-dev.service
Before=sysinit.target shutdown.target
Conflicts=shutdown.target

[Service]
# Yields (exit 0) when systemd's RuntimeWatchdog already owns the device, so
# on-failure (not always) avoids respawning on a clean yield.
Type=exec
ExecStart=/usr/lib/pocketforge/watchdog-keepalive.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=sysinit.target
