# rootfs-overlay/etc/systemd/system/pocketforge-splash-handoff.target
# pocketforge-splash-handoff.target — the takeover signal for the boot animator.
#
# The animator (pocketforge-boot-animator.service) has Conflicts= this target,
# so activating it stops the animator (systemd sends SIGTERM). The animator
# then clears fb0 to black and exits 0, leaving the panel clean for whatever
# owns the display next.
#
# A future MainUI / kiosk service wires itself in with:
#     [Unit]
#     Requires=pocketforge-splash-handoff.target
#     After=pocketforge-splash-handoff.target
# and the boot ordering hands the display off automatically. No sentinel file,
# no ad-hoc socket — one service edit is the entire integration.
#
# bd: tsp-3rd3.4

[Unit]
Description=PocketForge splash-to-UI handoff signal (takeover from boot animator)
Documentation=man:systemd.target(5)
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target
