# rootfs-overlay/etc/systemd/system/pocketforge-splash-handoff-default.service
# pocketforge-splash-handoff-default.service — interim safety net.
#
# Until a real MainUI/kiosk exists, nothing else activates
# pocketforge-splash-handoff.target, so the animator would run forever. This
# oneshot activates the target once multi-user.target has settled — the
# animator receives SIGTERM, clears fb0 to black, exits 0, and the panel
# lands in a defined state instead of a stuck animation loop.
#
# When a real UI service arrives it will bring its OWN dependency on
# pocketforge-splash-handoff.target (Requires=/After=), the handoff fires
# earlier, and this oneshot becomes a redundant no-op (activating an already-
# active target is a no-op). Delete this unit only after the real successor
# is landed in every variant.
#
# bd: tsp-3rd3.4

[Unit]
Description=PocketForge default splash-handoff (activate handoff.target when no UI is present)
Documentation=https://github.com/pocketforge-os/image/tree/main/apps/pocketforge-boot-animator
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
# Small settle so the login prompt has painted its own last state before we
# tear the animator down (avoids a visible 1-frame race on getty startup).
ExecStartPre=/bin/sleep 0.5
ExecStart=/bin/systemctl start pocketforge-splash-handoff.target
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
