# rootfs-overlay/etc/systemd/system/pocketforge-wifi-powersave.service
# =============================================================================
# pocketforge-wifi-powersave.service — Disable WiFi power-save on wlan0
# =============================================================================
# The XR819/XR829 (xradio) WiFi chip's 802.11 power-save implementation is
# unreliable: when the link goes quiet the chip sleeps, the AP stops receiving
# keepalives, ages the station out of its table, and then deauthenticates it
# with Reason 6 ("Class 2 frame from nonassociated STA"). The client immediately
# reassociates, producing the ~30s deauth/reassoc flap that prevents a stable
# DHCP lease from settling on the owned substrate.
#
# Disabling power-save is the default product policy (stability-first); it is a
# per-wdev nl80211 property that persists across association cycles, so a single
# application at interface bring-up is sufficient. The policy is overridable via
# a POWER_SAVE key in /boot/wifi.txt — see wifi-powersave.sh.
#
# bd: tsp-cv7.4.12
# =============================================================================
[Unit]
Description=PocketForge: apply WiFi power-save policy to wlan0 (default off — xradio flap mitigation)
# wlan0 is created when the xradio module triplet loads (modules-load.d). Order
# after the device node and after wpa_supplicant so the interface exists.
# bd tsp-mc9m.14.8: this unit is pulled in via the wlan0 .device unit's .wants/
# (see [Install]) rather than multi-user.target, so it is only in the boot
# transaction when wlan0 actually appears. That is why there is NO
# Wants=sys-subsystem-net-devices-wlan0.device here: a soft Wants= on an absent
# device still enqueues a start job that After= waits on, timing out at
# DefaultDeviceTimeoutSec (~90s) and stalling boot on a driver-less mainline
# kernel. Being device-.wants pulled avoids that entirely; After= below keeps the
# ordering (after the device is active + after wpa_supplicant) when wlan0 exists.
After=sys-subsystem-net-devices-wlan0.device wpa_supplicant@wlan0.service
# Skip cleanly on hardware/builds where wlan0 never appears.
ConditionPathExists=/sys/class/net/wlan0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/pocketforge/wifi-powersave.sh

[Install]
# Device-Wants-service (hotplug): WantedBy the wlan0 device, NOT multi-user.target
# — so an absent wlan0 never pulls this unit into the boot transaction (no ~90s
# device-timeout stall). build-rootfs.sh symlinks it into the matching
# sys-subsystem-net-devices-wlan0.device.wants/ directory.
WantedBy=sys-subsystem-net-devices-wlan0.device
